# docker-components/uwsgi/uwsgi-supervisor.conf
[program:uwsgi]
command = uwsgi --ini /var/www/code/docker-components/uwsgi/uwsgi.ini
